# src/etc/init.d/bell
#!/bin/sh /etc/rc.common
START=99
start() {
  bell off
  sleep 10
  bell listen > /tmp/bell.log &
}
